# .config/systemd/user/checkmail.timer
[Unit]
Description=Sync mail with isync
RefuseManualStart=no
RefuseManualStop=no

[Timer]
Persistent=false
OnBootSec=1min
OnUnitActiveSec=5min
Unit=checkmail.service

[Install]
WantedBy=timers.target
